# dpkg/info/aptitude.list
/.
/etc
/etc/cron.daily
/etc/cron.daily/aptitude
/etc/logrotate.d
/etc/logrotate.d/aptitude
/usr
/usr/bin
/usr/bin/aptitude-curses
/usr/share
/usr/share/bug
/usr/share/bug/aptitude
/usr/share/doc
/usr/share/doc/aptitude
/usr/share/doc/aptitude/NEWS.Debian.gz
/usr/share/doc/aptitude/changelog.Debian.gz
/usr/share/doc/aptitude/copyright
/usr/share/menu
/usr/share/menu/aptitude
/var
/var/lib
/var/lib/aptitude
